# services/dpi-tcpdump.service
[Unit]
Description=DPI Lab tcpdump ring buffer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/dpi-lab/scripts/run_tcpdump.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
